# provisioning/roles/ruby/files/fakes3_supervisor.conf
[program:fakes3]
command=/usr/local/bin/fakes3 -r /deploy/sites/fakes3 -p 4567
directory=/deploy/sites/fakes3
user=vagrant
autostart=true
autorestart=true
redirect_stderr=True