# provisioning/roles/ruby/files/mailcatcher_supervisor.conf
[program:mailcatcher]
command=/usr/local/bin/mailcatcher --smtp-port 10025 --http-ip 0.0.0.0 -f
directory=/deploy/sites
user=vagrant
autostart=true
autorestart=true
redirect_stderr=True