# provisioning/roles/thumbor/files/supervisor.conf
[program:thumbor]
command=/deploy/virtualenvs/thumbor/bin/thumbor -c thumbor.conf
directory=/deploy/sites/thumbor
user=vagrant
autostart=true
autorestart=true
redirect_stderr=True